# services/pipinstall.service
[Unit]
Description=pip updater for urlmon

[Service]
User=root
ExecStartPre=source bin/activate
ExecStart=/root/urlmon/bin/python3 -m pip install -r requirements.txt
WorkingDirectory=/root/urlmon

[Install]
WantedBy=multi-user.target